# syft/pkg/cataloger/debian/testdata/info/util-linux.conffiles
/etc/default/hwclock
/etc/init.d/hwclock.sh
/etc/pam.d/runuser
/etc/pam.d/runuser-l
/etc/pam.d/su
/etc/pam.d/su-l